# supervisor.conf
[program:shorturl]
command=/path/to/shorturl -redis redis://xxx
autostart=true
autorestart=true
startsecs=5
user=root
redirect_stderr = true
stdout_logfile = /var/log/supervisord/%(program)s.log